# code/bike_project/crontab
*/1 * * * * root python /bike_logs.py >> /data/cron.log
